# docker/s6/hermitshell-agent/run
#!/command/execlineb -P
/usr/local/bin/hermitshell-agent
