# systemd/hermitshell-agent.service
[Unit]
Description=HermitShell Router Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/hermitshell/hermitshell-agent
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hermitshell-agent
Environment=RUST_LOG=info
ProtectHome=yes
ProtectSystem=strict
ReadWritePaths=/data/hermitshell /run/hermitshell /tmp
PrivateTmp=yes
NoNewPrivileges=yes
PrivateDevices=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictSUIDSGID=yes
ProtectClock=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
RestrictRealtime=yes

[Install]
WantedBy=multi-user.target
