# hermitshell-agent/hermitshell-agent.service
[Unit]
Description=HermitShell Agent
After=network.target

[Service]
Type=simple
ExecStart=/opt/hermitshell/hermitshell-agent
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
